# auto-suspend-daemon.service
[Unit]
Description=Start auto-suspend in daemon mode

[Service]
Type=idle
ExecStart=/usr/bin/auto-suspend daemon

[Install]
WantedBy=multi-user.target
