# auto-suspend-clerk.service
[Unit]
Description=Run auto-suspend

[Service]
Type=oneshot
ExecStart=/usr/sbin/auto-suspend clerk
RemainAfterExit=no
